# projector-bridge.service
[Unit]
Description=Sony Projector IR-to-ADCP Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=trent
Group=input
WorkingDirectory=/home/trent/sony-projector-ip-control
ExecStartPre=+/usr/bin/ir-keytable -s rc0 -p sony
ExecStart=/home/trent/sony-projector-ip-control/.venv/bin/projector-bridge --log-level INFO
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/trent/sony-projector-ip-control /sys/class/rc/rc0/
PrivateTmp=true

[Install]
WantedBy=multi-user.target
